# _home/profile.d/00-env.sh
#!/usr/bin/env bash
# Environment variables for login shells
# This file is sourced by both zprofile and bash_profile

# If you use NVM, it's init would go here, but we use mise to manage node versions
